# roles/catalogue/templates/catalogue.service.j2
[Unit]
Description = Catalogue Service

[Service]
User=roboshop
Environment=MONGO=true
// highlight-start
Environment=MONGO_URL="mongodb://{{ MONGODB_HOST }}:27017/catalogue"
<<<<<<< HEAD
Environment=TOOL="ANSIBLE"
=======
>>>>>>> f083818b4add734bcc3977b2abd360034cb412f6
// highlight-end
ExecStart=/bin/node /app/server.js
SyslogIdentifier=catalogue

[Install]
WantedBy=multi-user.target