# roles/cart/templates/cart.service.j2
[Unit]
Description = Cart Service
[Service]
User=roboshop
// highlight-start
Environment=REDIS_HOST={{ REDIS_HOST }}
Environment=CATALOGUE_HOST={{ CATALOGUE_HOST }}
<<<<<<< HEAD
Environment=CATALOGUE_PORT=80
=======
Environment=CATALOGUE_PORT=8080
>>>>>>> f083818b4add734bcc3977b2abd360034cb412f6
// highlight-end
ExecStart=/bin/node /app/server.js
SyslogIdentifier=cart

[Install]
WantedBy=multi-user.target